# deploy/weekly-digest.timer
[Unit]
Description=WeeklyDigestOrbit - Weekly Timer (毎週月曜 AM 1:00)

[Timer]
OnCalendar=Mon *-*-* 01:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
